# etc/systemd/r3d-daemon.service
[Unit]
Description=rebelove.org command daemon

[Service]
ExecStart=/opt/rebelove.org/sbin/daemon.sh

[Install]
WantedBy=multi-user.target
